# birthing_suite/systemd/checkpoint.timer
[Unit]
Description=Daily checkpoint for Mem0 store

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
